# 99-ryzen-tdp.rules
 ACTION=="add", RUN+="/bin/sh -c 'chmod 0222 /sys/devices/platform/asus-nb-wmi/{ppt_fppt,ppt_pl1_spl,ppt_pl2_sppt}'"
